# systemd/recordara-sync.service
[Unit]
Description=Recordara sync

[Service]
ExecStart=/opt/recordara/do_sync.sh
